# roles/estacao/templates/cmc-boot.timer.j2
[Unit]
Description=Boot script

[Timer]
OnBootSec=2h
OnUnitActiveSec=6h

[Install]
WantedBy=timers.target
