# roles/estacao/templates/cmc-boot.timer.j2
[Unit]
Description=Boot script

[Timer]
OnBootSec=0min
OnUnitActiveSec=12h

[Install]
WantedBy=timers.target
